# docker/php/cron/crontabs/www-data
* * * * * /usr/local/bin/php /var/www/html/artisan auth:prune-expired-link-registration-tokens
* * * * * /usr/local/bin/php /var/www/html/artisan auth:prune-expired-link-login-tokens
